# cloudflare-ddns/rootfs/etc/services.d/cloudflare-ddns/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Add-on: Cloudflare DDNS
# Take down the service
# ==============================================================================

if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services
