# cloudflare-ddns/rootfs/etc/services.d/cloudflare-ddns/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Cloudflare DDNS
# Runs the Cloudflare DDNS service
# ==============================================================================

exec /run.sh
